# deploy/vaultic-api.service
# systemd service for the Vaultic FastAPI backend + Litestream replication.
#
# start.sh loads Litestream credentials from .env and wraps uvicorn with
# `litestream replicate -exec`, so replication starts before the app and
# stops cleanly when the app exits.
#
# One-time setup on server:
#   sudo cp deploy/vaultic-api.service /etc/systemd/system/vaultic-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable vaultic-api
#   sudo systemctl start vaultic-api
#
# NOTE: EnvironmentFile is intentionally omitted — systemd mangles multi-line
# PEM keys (Coinbase EC private key). load_dotenv() in api/main.py handles
# .env loading for the app. Litestream credentials are loaded by start.sh.

[Unit]
Description=Vaultic FastAPI Backend + Litestream Replication
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/vaultic
ExecStart=/bin/bash /home/ubuntu/vaultic/deploy/start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
